# provisioning/pipewire-system.service
[Unit]
Description=PipeWire audio server (system-wide, Venom appliance)
After=bluetooth.service

[Service]
Type=simple
User=venom
SupplementaryGroups=audio bluetooth
RuntimeDirectory=pipewire
RuntimeDirectoryPreserve=yes
Environment=XDG_RUNTIME_DIR=/run/pipewire
Environment=PIPEWIRE_RUNTIME_DIR=/run/pipewire
Environment=DBUS_SESSION_BUS_ADDRESS=unix:path=/run/dbus/system_bus_socket
ExecStart=/usr/bin/pipewire
Restart=always
RestartSec=3

[Install]
WantedBy=multi-user.target
